# integration/scripts/files/elements/guest-agent/install.d/guest-agent-source-install/docker-hostnic.socket
[Unit]
Description=docker hostnic driver

[Socket]
ListenStream=/run/docker/plugins/docker-hostnic.sock

[Install]
WantedBy=sockets.target