# integration/scripts/files/elements/guest-agent/install.d/guest-agent-source-install/guest-log-collection.timer
[Unit]
Description=OpenStack Trove Guest Log Collection Timer

[Timer]
OnActiveSec=10s
OnUnitInactiveSec=10s
Unit=guest-log-collection.service
Persistent=true

[Install]
WantedBy=multi-user.target